# conf/celery.service
[Unit]
Description=Celery Worker for Sanayicin
After=network.target redis-server.service

[Service]
Type=simple
User=sanayicin
Group=www-data
WorkingDirectory=/var/www/sanayicin/backend
Environment="PATH=/var/www/sanayicin/backend/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=main.settings_production"
ExecStart=/var/www/sanayicin/backend/venv/bin/celery -A main worker -l info --concurrency=2
Restart=always

[Install]
WantedBy=multi-user.target
